# macos/com.user.sync_all.plist
<!--
- com.user.sync_all.plist
- Copy to ~/Library/LaunchAgents to run while logged in and awake.
- Replace /Users/USERNAME with your actual home directory path.
-->

<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.configsync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>exec "$HOME/.config/scripts/sync_all_wrapper.sh"</string>
    </array>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>StandardOutPath</key>
    <string>/tmp/config_sync_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/config_sync_stderr.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
